# systemd/heyroya-beat.service
[Unit]
Description=HeyRoya Automation Celery Beat Scheduler
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/heyroya-automation
EnvironmentFile=/opt/heyroya-automation/.env
ExecStart=/opt/heyroya-automation/.venv/bin/celery -A app.workers.celery_app beat --loglevel=info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
